# zookeeper/zookeeper_supervisor.ini
[program:zookeeper]
user=nova
group=nova
command=/bin/bash /opt/ea/nova/zookeeper/bin/zookeeper.sh
priority=1
autostart=true
startretries=0
autorestart=false
exitcodes=0,1,2
stopsignal=TERM
stopwaitseconds=10
redirect_stderr=true
stdout_logfile=/var/log/zookeeper_supervisord.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
